# preload-input/15-ssh-host-keys/DEBIAN/postinst
#!/bin/sh
# DopamineX preload: ensure sshd host keys exist after openssh-server is installed.
#
# Why this exists:
#   openssh-server's own postinst runs ssh-keygen -A on FIRST install only.
#   Dopamine reinstalls openssh-server via 'dpkg -i' on every activation
#   (preinstalledDebs loop), and dpkg -i on an already-installed-same-version
#   package does NOT rerun postinst. So if the very first activation failed
#   midway (before openssh-server postinst could run keygen), all subsequent
#   activations silently skip key generation.
#
#   Symptom: ssh client sees "kex_exchange_identification: Connection closed
#   by remote host" -- sshd accepts TCP but the worker process exits before
#   sending its banner because it can't load any host key.
#
# What this does:
#   1. Generate each missing host key (rsa / ecdsa / ed25519) idempotently
#   2. Kickstart sshd via launchctl so a running instance picks up new keys
#
# Safe to re-run on every activation; pre-existing keys are not regenerated.

set -e

JB=/var/jb
SSH_DIR="$JB/etc/ssh"
KEYGEN="$JB/usr/bin/ssh-keygen"

# Sanity: openssh-client must have installed ssh-keygen before us
if [ ! -x "$KEYGEN" ]; then
    echo "DopamineX/15-ssh-host-keys: $KEYGEN not found; openssh-client missing?" >&2
    exit 0   # not fatal: skip gracefully so preload loop continues
fi

mkdir -p "$SSH_DIR"
chmod 0755 "$SSH_DIR"

# ─── Set default root password to "alpine" ─────────────────────
#
# !!!!!!!!!! SECURITY WARNING — KNOWN-PUBLIC PASSWORD !!!!!!!!!!
# This intentionally hard-codes the legacy jailbreak default password
# "alpine" so that fresh DopamineX installs are ssh-able as root without
# device-side console interaction. THIS REPO IS PUBLIC. ANY device with
# the resulting .tipa installed has root:alpine — anyone who knows the
# device's reachable address can ssh in as root.
#
# Acceptable only because this fork is for personal/dev use. If you
# clone this for someone else's device, change "alpine" below to a
# private secret OR rip out this entire block and rely on
# 'sudo passwd root' after first install.
# ───────────────────────────────────────────────────────────────
DEFAULT_PWD="alpine"
if command -v chpasswd >/dev/null 2>&1; then
    echo "DopamineX/15-ssh-host-keys: setting root password via chpasswd" >&2
    echo "root:${DEFAULT_PWD}" | chpasswd 2>/dev/null || true
elif command -v openssl >/dev/null 2>&1; then
    echo "DopamineX/15-ssh-host-keys: setting root password via master.passwd edit" >&2
    PWD_HASH=$(openssl passwd -6 "$DEFAULT_PWD" 2>/dev/null)
    if [ -n "$PWD_HASH" ]; then
        for pwfile in /var/jb/etc/master.passwd /var/jb/etc/shadow; do
            if [ -f "$pwfile" ]; then
                # 用 sed 替换 root 行的 password 列（第二列）
                # 用 | 作 sed 分隔符避免 hash 内含 / 冲突
                sed -i.dpxbak "s|^root:[^:]*:|root:${PWD_HASH}:|" "$pwfile" 2>/dev/null || true
            fi
        done
        # 如果 BSD-style: master.passwd 改完后需要 pwd_mkdb 重建 pwd.db/spwd.db
        if [ -f /var/jb/etc/master.passwd ] && command -v pwd_mkdb >/dev/null 2>&1; then
            pwd_mkdb -p /var/jb/etc/master.passwd 2>/dev/null || true
        fi
    fi
else
    echo "DopamineX/15-ssh-host-keys: WARN no chpasswd/openssl, skipping default password" >&2
fi

# ─── Patch sshd_config: iOS rootless 没 PAM ─────────────────────
# Default sshd_config ships UsePAM yes; on iOS rootless there is no
# /var/jb/etc/pam.d/sshd and no PAM .so modules, so worker dies with
# 'PAM: initialisation failed for "root"' the moment it tries to auth
# any user. Append rootless-safe override (idempotent: only writes if
# the marker is absent).
SSHD_CONF="$SSH_DIR/sshd_config"
MARKER="# DopamineX: rootless iOS auth overrides"
if [ -f "$SSHD_CONF" ] && ! grep -q "$MARKER" "$SSHD_CONF" 2>/dev/null; then
    echo "DopamineX/15-ssh-host-keys: patching $SSHD_CONF for rootless auth" >&2
    cat >> "$SSHD_CONF" <<EOF

$MARKER
UsePAM no
PasswordAuthentication yes
PermitRootLogin yes
KbdInteractiveAuthentication no
ChallengeResponseAuthentication no
EOF
fi

generated=0
for t in rsa ecdsa ed25519; do
    keyfile="$SSH_DIR/ssh_host_${t}_key"
    if [ ! -f "$keyfile" ]; then
        echo "DopamineX/15-ssh-host-keys: generating $keyfile" >&2
        "$KEYGEN" -q -t "$t" -N "" -f "$keyfile"
        chmod 0600 "$keyfile"
        chmod 0644 "${keyfile}.pub"
        generated=$((generated + 1))
    fi
done

# Bootstrap or restart sshd:
#   1. If launchd has the daemon registered  → kickstart (restart with new key)
#   2. If launchd has NOT registered it yet   → bootstrap (load plist from disk)
#
# Background: Dopamine's jbctl startup runs
#   launchctl bootstrap system /Library/LaunchDaemons
# but if that step never executes (e.g. first activation aborts midway),
# sshd's plist sits on disk unloaded -- ssh client sees Connection reset
# because nothing listens on the configured port. This makes the 15-host-keys
# package self-sufficient: it guarantees sshd is running before exiting postinst.

echo "DopamineX/15-ssh-host-keys: ensuring sshd is running" >&2

LAUNCHCTL="$JB/bin/launchctl"
[ -x "$LAUNCHCTL" ] || LAUNCHCTL="$(command -v launchctl || true)"

SSHD_LABEL="com.openssh.sshd"
SSHD_PLIST="/var/jb/Library/LaunchDaemons/${SSHD_LABEL}.plist"

# launchd 拒载 owner 不是 root 或 mode 过松的 plist (会报
# 'Bootstrap failed: 122: Path had bad ownership/permissions').
# dpkg 在某些路径下安装可能让属主变成 mobile 之类——主动 fix.
if [ -f "$SSHD_PLIST" ]; then
    chown root:wheel "$SSHD_PLIST" 2>/dev/null || true
    chmod 0644 "$SSHD_PLIST" 2>/dev/null || true
fi
# 上级目录链也得是 root + 0755；jb 安装可能漏修
for d in /var/jb/Library/LaunchDaemons /var/jb/Library /var/jb; do
    [ -d "$d" ] || continue
    chown root:wheel "$d" 2>/dev/null || true
    chmod 0755 "$d" 2>/dev/null || true
done

if [ -n "$LAUNCHCTL" ] && [ -f "$SSHD_PLIST" ]; then
    # 1. 是否已注册到 system domain?
    if "$LAUNCHCTL" print "system/${SSHD_LABEL}" >/dev/null 2>&1; then
        # 已注册 — 仅当新生成了 key 时才需重启
        if [ "$generated" -gt 0 ]; then
            echo "DopamineX/15-ssh-host-keys: kickstart sshd (existing daemon, ${generated} new key)" >&2
            "$LAUNCHCTL" kickstart -k "system/${SSHD_LABEL}" 2>&1 || true
        fi
    else
        # 未注册 — bootstrap 加载
        echo "DopamineX/15-ssh-host-keys: sshd not loaded; bootstrapping ${SSHD_PLIST}" >&2
        "$LAUNCHCTL" bootstrap system "$SSHD_PLIST" 2>&1 || \
            echo "DopamineX/15-ssh-host-keys: WARN bootstrap failed (non-fatal, continuing)" >&2
    fi
else
    echo "DopamineX/15-ssh-host-keys: WARN launchctl or sshd plist missing; skipping daemon mgmt" >&2
fi

if [ "$generated" -eq 0 ]; then
    echo "DopamineX/15-ssh-host-keys: all host keys already present" >&2
fi

exit 0
