# preload-input/40-toorless-roothide/Library/LaunchDaemons/com.dog.cat.toorless.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dog.cat.toorless</string>

    <!-- roothide: 第三方 /Library/LaunchDaemons 用【无前缀的 jbroot-based 路径】，
         roothide 的 launchd hook 在加载时解析到真实随机 jbroot。
         即把 rootless 的 /var/jb 前缀去掉即可。不要写 /var/jb，也不要写 @JBROOT@
         （@JBROOT@ 只对 Dopamine 自带的 /basebin/LaunchDaemons 生效，第三方不会被替换）。
         实证：roothide bootstrap 自带的 /Library/LaunchDaemons/com.apple.atrun.plist
         用的就是 "/usr/libexec/atrun"。 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/toorless</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DISABLE_TWEAKS</key>
        <string>1</string>
    </dict>

    <key>UserName</key>
    <string>root</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
