# preload-input/16-ssh-port-roothide/DEBIAN/postinst
#!/bin/sh
# DopamineX roothide preload: 把 sshd 的监听端口从 22/2222 改成 18888。
#
# 为什么需要：
#   roothide 的 openssh-server 用 launchd inetd 模式，监听端口由
#   /Library/LaunchDaemons/com.openssh.sshd.plist 的 Sockets 决定（默认两个
#   listener：SockServiceName=ssh→22 和 =2222），**不是** sshd_config 的 Port。
#   RootHide Manager 的“SSH Server has been installed”警告是 TCP 探测
#   127.0.0.1:22 和 :2222（反汇编确认，非查文件）——只要这两个口在听就报警。
#   把 listener 收敛成单个 18888，既绕过该警告，SSH 仍可用（连接用 -p 18888）。
#
# 它如何在【首次越狱】就直接起在 18888（关键时序）：
#   roothide 激活流程 = finalizeBootstrap（按序 dpkg -i 所有预装包，含 openssh-server
#   与本包）→ 紧接着 rebootUserspace（用户态重启）。dpkg -i 不会跑 openssh 的
#   extrainst_(`launchctl load`)，且 DOJailbreaker 里那句显式
#   `launchctl bootstrap system /Library/LaunchDaemons` 是注释掉的（注释原文：
#   “只有不立即 userspace reboot 时才需要”）。也就是说：
#     ① 本包 postinst 在 reboot 之前就把 plist 落成 18888；
#     ② 用户态重启后，launchdhook/daemon_hook 从磁盘的 /Library/LaunchDaemons 重新
#        加载，读到的就是 18888 的 plist → sshd 直接在 18888 起来。
#   所以前提只有一个：【写出去的 plist 必须有效】。故这里【无条件整份写一份已知良好的
#   XML plist】，不依赖 plutil、不在二进制 plist 上做易错的原地编辑。
#
# 安全红线：【绝不 launchctl bootout 活动 daemon】。曾经用 bootout+bootstrap，bootout
#   把 sshd 干掉、bootstrap 又因 daemon_hook 已登记同名服务(race)/权限失败没拉起，
#   导致 sshd 全端口消失把人锁死。现在只写文件 + 顶多一次纯 bootstrap，永不卸载。
#
# 注意：不用 set -e（dpkg postinst 环境 PATH/locale 不可控，单行失败不应中断后续）。
#       路径用无前缀绝对路径，libroot 会把 /Library /usr 重定向到随机 jbroot。

PORT=18888
PLIST="/Library/LaunchDaemons/com.openssh.sshd.plist"
LABEL="com.openssh.sshd"

log() { echo "DopamineX/16-ssh-port: $*" >&2; }

if [ ! -f "$PLIST" ]; then
    log "WARN $PLIST not found (openssh-server 未安装？) — 跳过"
    exit 0
fi

# ── 无条件整份重写为已知良好的 XML（单 listener=18888）──────────────────
# 内容对应 roothide procursus 的 openssh-server.deb（inetd wrapper 路径稳定）；
# 唯一改动 = Sockets 收敛成单个 18888 listener（去掉 22 与 2222）。
# launchd 同时接受 XML / 二进制 plist，XML 最省事且确定有效。
log "writing plist → single listener on $PORT"
cat > "$PLIST" <<PLIST_EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.openssh.sshd</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/bin/sh</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>/usr/libexec/sshd-keygen-wrapper</string>
		<string>-i</string>
	</array>
	<key>SessionCreate</key>
	<true/>
	<key>Sockets</key>
	<dict>
		<key>SSHListener</key>
		<dict>
			<key>SockServiceName</key>
			<string>${PORT}</string>
		</dict>
	</dict>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
	<key>inetdCompatibility</key>
	<dict>
		<key>Wait</key>
		<false/>
	</dict>
	<key>ExecuteAllowed</key>
	<true/>
</dict>
</plist>
PLIST_EOF

# 修属主/权限（launchd 拒载 owner 非 root 或 mode 过松的 plist：'bad ownership/permissions'）
chown root:wheel "$PLIST" 2>/dev/null || true
chmod 0644 "$PLIST" 2>/dev/null || true

# ── 顶多一次【纯 bootstrap】，绝不 bootout ──────────────────────────────
# 正常首次越狱：随后的 userspace reboot 会由 daemon_hook 从这份 18888 plist 加载，
# 所以这里的 bootstrap 通常是冗余的（sshd 还没起 → 这次拉起在 18888；已起 → no-op）。
# 不做任何卸载/重启动作，保证任何路径下都不会把 sshd 弄没 → 永不锁死。
LAUNCHCTL=""
for c in /usr/bin/launchctl /bin/launchctl "$(command -v launchctl 2>/dev/null)"; do
    [ -n "$c" ] && [ -x "$c" ] && { LAUNCHCTL="$c"; break; }
done
if [ -n "$LAUNCHCTL" ]; then
    "$LAUNCHCTL" bootstrap system "$PLIST" 2>/dev/null || true
fi

log "done — plist=$PORT；激活末尾 userspace reboot 后 sshd 直接在 $PORT (连接: ssh -p $PORT root@设备)"
exit 0
